# rootfs/etc/services.d/nslcd/run
#!/usr/bin/with-contenv sh

nslcd --version

exec nslcd --nofork