# rootfs/etc/services.d/nslcd/run
#!/usr/bin/with-contenv sh

nslcd --version

OPTS=""
if [ -z "${DEBUG}" ]; then
    OPTS="--debug"
fi

exec nslcd --nofork "${OPTS}"
