# rootfs/etc/services.d/sshd/run
#!/usr/bin/with-contenv sh

OPTS=""
if [ -z "${DEBUG}" ]; then
    OPTS="-ddd"
fi

exec /usr/sbin/sshd -i -p 22 "${OPTS}"
