# rootfs/etc/services.d/sshd/run
#!/usr/bin/with-contenv sh

if [ ! -z "${DEBUG}" ]; then
    OPTS="${OPTS:-} -ddd"
fi

exec /usr/sbin/sshd -D -p 22 ${OPTS:-}
