# rootfs/etc/services.d/smb/run
#!/usr/bin/with-contenv sh

samba --version

OPTS=""
if [ -z "${DEBUG}" ]; then
    OPTS="--debuglevel=5"
fi

exec smbd --foreground --no-process-group --log-stdout "${OPTS}"
