# rootfs/etc/services.d/smb/run
#!/usr/bin/with-contenv sh

samba --version

exec smbd --foreground --no-process-group --log-stdout --debuglevel=5
