# rootfs/etc/services.d/smb/run
#!/usr/bin/with-contenv sh

samba --version

exec smbd -i
