# rootfs/etc/services.d/smb/run
#!/usr/bin/with-contenv sh

samba --version

exec samba -i --configfile=${SAMBA_CONFIG}
